# package/base-files/files/etc/init.d/sysctl
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=11
start() {

	[ -f "/bin/noeat" ] && {
		/bin/noeat
	}
	
	[ -f /etc/sysctl.conf ] && sysctl -p -e >&-
}
